# supervisor/supervisord-gin.conf
[program:golang-catalyst]
command = /usr/local/go/src/catalyst.Go/main
stopsignal=QUIT
stopasgroup=true
killasgroup=true
autostart=true
autorestart=true
redirect_stderr=true
